# systemd/openclaw-coordinator.service
[Unit]
Description=OpenClaw AgentCoordinator Slack broker
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=-/etc/openclaw/coordinator.env
ExecStart=/usr/local/bin/openclaw-coordinatorctl --config /etc/openclaw/coordinator-config.json serve-slack
Restart=always
RestartSec=5
NoNewPrivileges=yes
PrivateTmp=yes
UMask=0077
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/openclaw
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
PrivateDevices=yes

[Install]
WantedBy=multi-user.target
